# deploy/traditional/systemd/xboard-horizon.service
[Unit]
Description=Xboard Laravel Horizon queue workers
After=network.target redis-server.service mysql.service
Wants=redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/xboard/current
Environment=APP_ENV=production
ExecStart=/usr/bin/php artisan horizon
ExecStop=/usr/bin/php artisan horizon:terminate
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
